# db/crontab
# crontab do servidor do cliente
* * * * * /var/www/scripts/monitor.sh dba zoilo367
#* * * * * rm /var/www/html/teste2.teste
